# scripts/systemd/operator-auth-monitor.service
[Unit]
Description=operator Auth Expiry Monitor
After=network.target

[Service]
Type=oneshot
# setup-auth-system.sh renders this placeholder to the checked-out script path.
ExecStart=@
# Configure notification channels via environment
Environment=WARN_HOURS=2
# Environment=NOTIFY_PHONE=+1234567890
# Environment=NOTIFY_NTFY=operator-alerts

[Install]
WantedBy=default.target
